# n8n-monitor.service
[Unit]
Description=n8n Monitor and Backup Service
After=network.target

[Service]
Type=simple
User=wadmin
WorkingDirectory=/home/wadmin/n8n/n8nbk_python
ExecStart=/home/wadmin/n8n/n8nbk_python/venv/bin/python /home/wadmin/n8n/n8nbk_python/app.py
Restart=always
RestartSec=10

# 環境變數
Environment="PYTHONUNBUFFERED=1"

# 日誌設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=n8n-monitor

[Install]
WantedBy=multi-user.target
